# cron.d/crococheck
*/5 * * * * root /etc/crocodile/checker.sh
