# cron.d/crocodile
*/2 * * * * root /etc/crocodile/sender.py
