# cron.d/crococheck
*/5 * * * * root /etc/crocodile/checker.sh > /dev/null 2>&1
