# cron.d/crococheck
*/5 * * * * root /etc/crocodile/checker.sh > /var/log/supervisor/cron.log 2>&1
